# deployment/medexjob-backend.service
# AI assisted development
# Systemd service file for MedExJob Backend
# Place this file at: /etc/systemd/system/medexjob-backend.service

[Unit]
Description=MedExJob Backend Service
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=medexjob
Group=medexjob
WorkingDirectory=/opt/medexjob/backend

# Environment variables
Environment="SPRING_PROFILES_ACTIVE=prod"
# FTP Configuration for Hostinger
Environment="FILE_USE_FTP=true"
Environment="FILE_BASE_URL=https://medexjob.com"
Environment="FTP_HOST=82.180.143.101"
Environment="FTP_PORT=21"
Environment="FTP_USERNAME=u284488379"
Environment="FTP_PASSWORD=Shivam@3650"
Environment="FTP_REMOTE_DIR=uploads"
# Uncomment and set these if needed:
# Environment="JWT_SECRET=your_secret_here"
# Environment="DB_PASSWORD=your_db_password"

# Java options
ExecStart=/usr/bin/java -Xms512m -Xmx1024m -jar -Dspring.profiles.active=prod /opt/medexjob/backend/target/medexjob-backend-1.0.0.jar

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medexjob-backend

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

